# systemd/plugwise-pi.service
[Unit]
Description=Plugwise Pi Data Collection Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/plugwise_pi
Environment=PYTHONPATH=/home/pi/plugwise_pi
ExecStart=/home/pi/plugwise_pi/venv/bin/python -m plugwise_pi.collector
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/pi/plugwise_pi/data /home/pi/plugwise_pi/logs

[Install]
WantedBy=multi-user.target 