# systemd/plugwise-pi-api.service
[Unit]
Description=Plugwise Pi API Server
After=network.target
Wants=network.target
Requires=plugwise-pi.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/plugwise_pi
Environment=PYTHONPATH=/home/pi/plugwise_pi
ExecStart=/home/pi/plugwise_pi/venv/bin/python -m plugwise_pi.api
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/pi/plugwise_pi/data /home/pi/plugwise_pi/logs

[Install]
WantedBy=multi-user.target 